# systemd/claude-anywhere.service
[Unit]
Description=Claude Anywhere — All Bridges
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/claude-anywhere
ExecStart=/usr/bin/node start-all.mjs
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/claude-anywhere/.env

[Install]
WantedBy=multi-user.target
